# msic/systemd/deepin-face.service
[Unit]
Description=Deepin Face Authenticate Driver 

[Service]
User=root
Type=dbus
BusName=com.deepin.face
ExecStart=/usr/libexec/deepin-face

