# msic/systemd/deepin-face.service
[Unit]
Description=Deepin Face Authenticate Driver 

[Service]
User=root
Type=dbus
BusName=org.deepin.dde.Face1
ExecStart=/usr/libexec/deepin-face

